# infra/systemd/vici2-backup-mysql.service
[Unit]
Description=vici2 MySQL daily backup
Documentation=file:spec/runbooks/restore.md
After=network.target

[Service]
Type=oneshot
User=vici2-backup
EnvironmentFile=/etc/vici2/backup.env
ExecStart=/opt/vici2/scripts/backup/mysql.sh --env prod --archive-class daily
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
