# infra/systemd/vici2-backup-valkey-monthly.timer
[Unit]
Description=vici2 Valkey monthly backup timer
Documentation=file:spec/runbooks/restore.md

[Timer]
OnCalendar=*-*-01 02:45:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
