# infra/systemd/vici2-backup-valkey.timer
[Unit]
Description=vici2 Valkey daily backup timer
Documentation=file:spec/runbooks/restore.md

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
